# home/sovol/usbmount/makerbase-automount@.service
[Service]
Type=forking
GuessMainPID=no
ExecStart=/usr/bin/makerbase-automount %I
